# scripts/com.openclaw.sync-attachments.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.openclaw.sync-attachments</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/maria/.openclaw/workspace-main/scripts/sync-attachments.sh</string>
  </array>

  <key>WatchPaths</key>
  <array>
    <string>/Users/maria/Library/Messages/Attachments</string>
  </array>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/maria/.openclaw/workspace-main/inbox/.launchd.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/maria/.openclaw/workspace-main/inbox/.launchd-err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/maria</string>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
</dict>
</plist>
